# SSR/kcp.conf
[program:kcptun]
priority = 150
command = /opt/kcptun -c /opt/kcptun.json
